# modified_ac200/usr/lib/opkg/info/pingcheck.list
/etc/init.d/pingcheck
/etc/pingcheck/online.d/example-script.sh
/etc/pingcheck/offline.d/example-script.sh
/usr/sbin/pingcheck
/etc/config/pingcheck
