# squashfs-root/usr/lib/opkg/info/base-files.list
/etc/init.d/remote
/etc/defconfig/cf-wr305n/led_mode
/etc/defconfig/cf-e500n/version
/etc/defconfig/cf-wr305n/workmode
/etc/defconfig/cf-ac200/network
/etc/defconfig/cf-wr753ac/mbox
/etc/banner.failsafe
/etc/init.d/wl2tpaddroute
/etc/resolv.conf
/etc/board.d/01_leds
/etc/init.d/done
/etc/defconfig/cf-e330n/uhttpd
/bin/config_generate
/var
/etc/defconfig/cf-wr380ac/network
/etc/hotplug.d/block/30-usbmount
/lib/ramips.sh
/etc/init.d/system
/etc/preinit
/lib/dnsmasq_queue.sh
/etc/defconfig/cf-wr625ac/version
/bin/cluster_led_sh
/lib/functions/uci-defaults-ap.sh
/etc/defconfig/cf-ac100/version
/etc/defconfig/cf-e510n/uhttpd
/etc/defconfig/cf-wr755ac/led_mode
/lib/upgrade/platform.sh
/etc/rc.button/change_mode0
/etc/rc.button/change_mode1
/etc/shells
/etc/defconfig/cf-wr755ac/workmode
/lib/preinit/10_sysinfo
/etc/defconfig/cf-ac101/system
/lib/preinit/70_initramfs_test
/etc/config/led
/bin/repeater_status_ledcrtl
/etc/protocols
/etc/defconfig/cf-e320n/system
/etc/uci-defaults/12_network-generate-ula
/etc/config/remote
/etc/init.d/pcie_down_reboot
/etc/defconfig/cf-wr627n/network
/lib/preinit/10_indicate_preinit
/lib/preinit/03_preinit_do_ramips.sh
/etc/key
/etc/defconfig/cf-wr754ac/version
/etc/defconfig/cf-wr620n/network
/etc/defconfig/cf-ac200/system
/etc/defconfig/cf-wr305n/mbox
/etc/defconfig/cf-wr620n/uhttpd
/etc/defconfig/cf-e500n/system
/etc/defconfig/cf-ac101/led_mode
/lib/functions/uci-defaults-new.sh
/etc/defconfig/cf-ac101/workmode
/etc/defconfig/cf-wr754ac/mbox
/etc/defconfig/cf-wr618ac/led_mode
/etc/defconfig/cf-wr620n/led_mode
/sbin/reboot.sh
/etc/defconfig/cf-wr625ac/system
/etc/defconfig/cf-e330n/version
/etc/config/system
/etc/defconfig/cf-wr754ac/system
/etc/defconfig/cf-wr618ac/workmode
/etc/defconfig/cf-wr620n/workmode
/etc/defconfig/cf-wr800n/network
/usr/sbin/smp.sh
/etc/config/common_ap
/etc/defconfig/cf-wr800n/uhttpd
/etc/defconfig/cf-ac200/led_mode
/bin/generaterule.sh
/sbin/wps
/etc/defconfig/cf-wr618ac/network
/etc/fstab
/etc/sysctl.conf
/etc/defconfig/cf-ac200/workmode
/sbin/sysupgrade
/etc/defconfig/cf-wr380ac/uhttpd
/etc/config/mbox
/lib/preinit/99_10_run_init
/etc/defconfig/cf-wr800n/led_mode
/etc/defconfig/cf-e510n/version
/etc/hotplug.d/iface/11-customqos
/sbin/led.sh
/etc/defconfig/cf-wr800n/workmode
/lib/preinit/50_indicate_regular_preinit
/etc/config/vlan
/etc/config/login
/etc/shadow
/etc/rc.common
/etc/config/system_ac
/etc/defconfig/cf-wr305n/system
/etc/defconfig/cf-wr305n/version
/etc/uci-defaults/10_migrate-shadow
/etc/defconfig/cf-wr627n/system
/etc/defconfig/cf-ac101/version
/etc/config/system_ap
/etc/defconfig/cf-e320n/led_mode
/lib/preinit/10_indicate_failsafe
/init
/etc/defconfig/cf-e320n/workmode
/etc/board.d/02_network
/etc/config/probe
/etc/defconfig/cf-wr618ac/uhttpd
/etc/defconfig/cf-ac100/mbox
/etc/init.d/watchdog
/etc/hotplug.d/iface/01-fix-vpn-default-route
/etc/defconfig/cf-ac200/version
/sbin/mtk_wps
/etc/defconfig/cf-wr625ac/mbox
/etc/defconfig/cf-e500n/led_mode
/etc/defconfig/cf-wr753ac/led_mode
/usr/sbin/setsmp.sh
/lib/preinit/80_mount_root
/etc/config/cluster
/etc/defconfig/cf-e330n/system
/etc/defconfig/cf-wr627n/led_mode
/etc/defconfig/cf-e500n/workmode
/etc/defconfig/cf-wr753ac/workmode
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/defconfig/cf-wr627n/workmode
/etc/defconfig/cf-e320n/mbox
/bin/mount_disk
/etc/defconfig/cf-wr620n/mbox
/etc/defconfig/cf-wr380ac/version
/etc/init.d/led
/etc/init.d/usb_start
/etc/defconfig/cf-wr753ac/network
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/bin/config_generate_rp
/etc/init.d/configs
/bin/config_generate_rt
/etc/defconfig/cf-e510n/system
/etc/defconfig/cf-wr755ac/version
/lib/preinit/04_handle_checksumming
/etc/defconfig/cf-e500n/mbox
/etc/defconfig/cf-wr800n/mbox
/bin/button_change_mode
/etc/defconfig/cf-ac100/uhttpd
/etc/rc.button/wps
/bin/calcip
/etc/rc.button/failsafe
/etc/defconfig/cf-e320n/network
/etc/services
/etc/defconfig/cf-ac101/mbox
/lib/functions/system.sh
/etc/defconfig/cf-wr627n/version
/lib/functions/uci-defaults-4g.sh
/etc/defconfig/cf-wr625ac/led_mode
/bin/board_detect
/etc/defconfig/cf-wr753ac/uhttpd
/lib/upgrade/common.sh
/etc/defconfig/cf-wr625ac/workmode
/etc/defconfig/cf-e500n/network
/etc/defconfig/cf-wr620n/system
/etc/defconfig/cf-wr620n/version
/etc/rc.button/power
/etc/sysupgrade.conf
/etc/hotplug.d/firmware/10-rt2x00-eeprom
/etc/defconfig/cf-ac200/mbox
/etc/rc.button/rfkill
/lib/preinit/07_set_preinit_iface_ramips
/lib/upgrade/keep.d/base-files
/etc/init.d/customqos
/lib/preinit/40_run_failsafe_hook
/etc/defconfig/cf-wr800n/system
/etc/group
/etc/defconfig/cf-wr800n/version
/etc/rc.button/reset
/lib/preinit/30_failsafe_wait
/etc/defconfig/cf-e330n/mbox
/etc/defconfig/cf-wr380ac/system
/bin/led_usr
/etc/defconfig/cf-wr618ac/version
/sbin/firstboot
/etc/defconfig/cf-wr625ac/network
/etc/uci-defaults/09_fix-seama-header
/sbin/timing-reboot
/etc/defconfig/cf-ac100/network
/etc/defconfig/cf-e330n/led_mode
/bin/login.sh
/etc/hotplug.d/net/00-sysctl
/etc/defconfig/cf-e330n/workmode
/etc/hosts
/etc/defconfig/cf-e510n/mbox
/lib/functions.sh
/etc/init.d/arp_set
/etc/defconfig/cf-ac101/uhttpd
/etc/uci-defaults/11_migrate-sysctl
/lib/functions/leds.sh
/etc/defconfig/cf-e320n/uhttpd
/etc/defconfig/cf-e510n/led_mode
/sbin/remote_control
/etc/defconfig/cf-wr754ac/led_mode
/etc/defconfig/cf-wr627n/mbox
/rom/note
/etc/defconfig/cf-e510n/workmode
/etc/defconfig/cf-wr618ac/system
/etc/inittab
/etc/config/common
/etc/defconfig/cf-wr754ac/workmode
/etc/defconfig/cf-wr754ac/network
/sbin/wifi
/etc/defconfig/cf-ac200/uhttpd
/lib/wifi/ralink_common.sh
/etc/profile
/lib/network_queue.sh
/etc/defconfig/cf-e500n/uhttpd
/etc/defconfig/cf-wr618ac/mbox
/etc/init.d/wpptpaddroute
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/etc/defconfig/cf-e330n/network
/etc/TZ
/etc/mcc.json
/etc/defconfig/cf-wr625ac/uhttpd
/etc/init.d/umount
/lib/functions/service.sh
/etc/defconfig/cf-wr754ac/uhttpd
/etc/config/arp
/lib/upgrade/keep.d/base-files-essential
/etc/defconfig/cf-wr753ac/version
/etc/passwd
/lib/preinit/02_default_set_state
/etc/defconfig/cf-wr380ac/led_mode
/etc/hotplug.d/usb/10-motion
/etc/mtab
/etc/init.d/boot
/etc/defconfig/cf-wr380ac/workmode
/etc/defconfig/cf-e510n/network
/bin/button_light
/etc/defconfig/cf-ac100/system
/etc/defconfig/cf-e320n/version
/etc/defconfig/cf-ac100/led_mode
/lib/functions/preinit.sh
/bin/config_generate_ap
/bin/ipcalc.sh
/etc/defconfig/cf-wr305n/network
/etc/defconfig/cf-ac100/workmode
/etc/defconfig/cf-ac101/network
/etc/diag.sh
/etc/defconfig/cf-wr305n/uhttpd
/etc/defconfig/cf-wr627n/uhttpd
/etc/defconfig/cf-wr753ac/system
/etc/defconfig/cf-wr380ac/mbox
